# rclone-systemd-sync.timer
[Unit]
Description=Run sync script every 30 minutes and at boot

[Timer]
OnBootSec=1min
OnUnitActiveSec=30min
Unit=rclone-systemd-sync.service

[Install]
WantedBy=timers.target
